# debian/wold-server.service
[Unit]
Description=Wake On LAN Dashboard

[Service]

[Install]
WantedBy=multi-user.target